# php-8.1-fpm/defaults.d/apache2.conf
[program:fpm]
command=/usr/local/sbin/php-fpm