# nginx/defaults.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
redirect_stderr=true
redirect_stdout=true
