# nginx/defaults.d/php-fpm.conf
[program:fpm]
command=/usr/local/sbin/php-fpm -F -y /usr/local/etc/php-fpm.d/www.conf
logfile=/dev/stdout
logfile_maxbytes=0
