# nginx/defaults.d/php-fpm.conf
[program:fpm]
command=/usr/local/sbin/php-fpm -F -y /usr/local/etc/php-fpm.d/www.conf
redirect_stderr=true
redirect_stdout=true
