# nginx/defaults.d/php-fpm.conf
[program:fpm]
command=/usr/local/sbin/php-fpm
logfile=/dev/stdout
logfile_maxbytes=0