# apache/defaults.d/apache2.conf
[program:apache2]
command=/usr/local/bin/apache2-foreground
redirect_stderr=true
redirect_stdout=true
